# dockerized/scheduler/crontab
*/10 * * * * docker restart funda-scrapper-container
# extra line to make this a valid crontab.