# dockerized/scheduler/crontab
*/10 * * * * docker restart funda_scrapper-container
# extra line to make this a valid crontab.