# dockerized/scheduler/crontab
0 8 * * * docker restart scrapper-scripts-container && docker restart scapper-database-container
# extra line to make this a valid crontab.